# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Stop and disable old service if it exists
        if [ -f /etc/init.d/downtime ]; then
            if systemctl is-active --quiet downtime 2>/dev/null; then
                systemctl stop downtime || true
            fi
            if systemctl is-enabled --quiet downtime 2>/dev/null; then
                systemctl disable downtime || true
            fi
            # Remove old init.d script
            rm -f /etc/init.d/downtime
        fi

        # Enable and start new service
        systemctl daemon-reload
        systemctl enable checkmk-downtime.service
        systemctl start checkmk-downtime.service
        ;;
esac

#DEBHELPER#
